# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:create-branch]
command=/golang/create-branch/serve.sh /create-branch-conf

